# systemd/jetson-vision-desk.service
[Unit]
Description=Jetson Nano Vision Desk — HUD + vision/reflex loop on the DWIN screen
# the HUD needs an X session; on a stock JetPack desktop image that's the
# graphical target's :0. If you boot to console, run the loop with --demo instead.
After=graphical.target nv-l4t-bootloader-config.service
Wants=graphical.target

[Service]
Type=simple
User=jetson
# the desktop session owns DISPLAY :0. XAUTHORITY depends on the display manager:
#   gdm3   -> /run/user/1000/gdm/Xauthority   (JetPack 4.6 default; 1000 = `id -u jetson`)
#   lightdm-> /home/jetson/.Xauthority  (or /var/run/lightdm/...)
# check `echo $XAUTHORITY` inside the desktop session if it won't connect.
Environment=DISPLAY=:0
Environment=XAUTHORITY=/run/user/1000/gdm/Xauthority
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
# preload libgomp before CUDA libs — otherwise on aarch64 cv2's first `import cv2` after
# jetson-inference is loaded throws "libgomp.so.1: cannot allocate memory in static TLS block"
Environment=LD_PRELOAD=/usr/lib/aarch64-linux-gnu/libgomp.so.1
WorkingDirectory=/home/jetson/jetson-vision-desk
# wait until the WM (compiz on JetPack 4.6's Unity desktop) is actually running before
# launching the HUD — graphical.target reaches before Compiz is up, and starting the HUD
# pre-Compiz makes the WM later decorate our window and ignore the fullscreen hint.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do pgrep -x compiz >/dev/null 2>&1 && break; sleep 1; done; sleep 2'
# kiosk: silence notification banners so they don't pop over the fullscreen HUD
# (leading '-' = don't fail the unit if gsettings/dbus isn't there)
ExecStartPre=-/usr/bin/gsettings set org.gnome.desktop.notifications show-banners false
ExecStart=/usr/bin/python3 -m src.main --hud
Restart=on-failure
RestartSec=5
TimeoutStartSec=90

[Install]
WantedBy=graphical.target
